# deploy/feedling-backend.service
[Unit]
Description=Feedling Flask Backend
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/feedling-mcp/backend
EnvironmentFile=/home/ubuntu/feedling.env
ExecStart=/home/ubuntu/feedling-venv/bin/gunicorn --config gunicorn_conf.py -w 1 -k asgi.worker.FeedlingUvicornWorker --timeout 120 -b 0.0.0.0:5001 asgi_app:app
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
